# configs/supervisord.conf
; Supervisord Configuration 
; Supervisor is a process manager used to manage Plausible & Clickhouse processes
; This file is generated by Cloudron. Do not modify.
; For more information on the config file, please see:
; http://supervisord.org/configuration.html

[unix_http_server]
file=/run/supervisord/supervisor.sock   ; the path to the socket file
chown=cloudron:cloudron       ; socket file uid:gid owner
;chmod=0700                 ; socket file mode (default 0700)
;username=user              ; default is no username (open server)
;password=123               ; default is no password (open server)

[supervisord]
pidfile=/run/supervisord/supervisord.pid ; supervisord pidfile; default supervisord.pid
loglevel=info                ; log level; default info; others: debug,warn,trace
logfile=/run/supervisord/supervisord.log ; main log file; default $CWD/supervisord.log
; Cloudron takes care of log rotation automatically. Hence:
logfile_maxbytes=0           ; max main logfile bytes b4 rotation; Set to unlimited.
logfile_backups=0            ; # of main logfile backups; 0 means none.
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
; We run the supervisord daemon as root, because we are dropping privileges later in
; the [program:plausible] and [program:clickbase] sections using the user directive.
; This is not a security issue as the actual plausible and clickbase daemons are run as
; less privileged accounts
user=root                    ; setuid to this UNIX account at startup; recommended if root
childlogdir=/tmp             ; 'AUTO' child log dir, default $TEMP
; environment=KEY="value"    ; key value pairs to add to environment. We do not need to
; set this because supervisor subprocesses already inherit the environment (i.e. 
; the ENV variables) of the shell used to start supervisord. See:
; http://supervisord.org/subprocess.html#subprocess-environment

; The rpcinterface:supervisor section must remain in the config file for
; RPC (supervisorctl/web interface) to work.  Additional interfaces may be
; added by defining them in separate [rpcinterface:x] sections.
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; The supervisorctl section configures how supervisorctl will connect to
; supervisord.  configure it match the settings in either the unix_http_server
; or inet_http_server section.
[supervisorctl]
serverurl=unix:///run/supervisord/supervisor.sock ; use a unix:// URL  for a unix socket
;username=chris              ; should be same as in [*_http_server] if set
;password=123                ; should be same as in [*_http_server] if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

; The sample program section below shows all possible program subsection values.
; Create one or more 'real' program: sections to be able to control them under
; supervisor.

[program:clickbase]
command=/usr/bin/clickhouse-server --config-file /app/data/clickhouse-config.xml --pid-file /run/clickhouse/clickhouse-server.pid
user=clickhouse
priority=1
autostart=true
autorestart=true
startsecs=5
; We will have Docker manage the log outputs. See:
; https://veithen.io/2015/01/08/supervisord-redirecting-stdout.html
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:plausible]
command=/app/code/plausible/bin/plausible start
user=cloudron
priority=2
autostart=true
autorestart=true
startsecs=5
; We will have Docker manage the log outputs. See:
; https://veithen.io/2015/01/08/supervisord-redirecting-stdout.html
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0